# DEBIAN/postinst
#!/bin/bash

systemctl enable trains
systemctl restart trains