# DEBIAN/postinst
#!/bin/bash

systemctl enable trains nginx
systemctl restart trains nginx